# home/dot_config/systemd/user/zsh-history-backup.service
[Unit]
Description=Back up zsh history (dedup'd union archive + rotated snapshots)

[Service]
Type=oneshot
ExecStart=%h/.local/bin/zsh-history-backup
